# files/cron-supervisor.conf
[program:cron]
priority=10
command = cron -f -L 15
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true
stopsignal=QUIT
