# files/noVNC-supervisor.conf
[program:novnc]
priority=25
directory=/noVNC
command=/noVNC/utils/launch.sh
user=nobody
autostart=true
autorestart=true
stopsignal=QUIT
