# files/rsyslog-supervisor.conf
[program:rsyslog]
priority=10
command = rsyslogd -n
stopwaitsecs = 5
stopsignal=QUIT
stdout_logfile=/tmm/rlogs/supervisor/rsyslog.log
stderr_logfile=/tmm/rlogs/supervisor/rsyslog.err
user=root
