# files/rsyslog-supervisor.conf
[program:rsyslog]
priority=10
command = rsyslogd -n -c3
stopwaitsecs = 5
stopsignal=QUIT
