# files/tmm-supervisor.conf
[program:tmm]
priority=30
directory=/tmm
command=/tmm/tinyMediaManager.sh
user=nobody
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1",HOME="/nobody"

