# files/x11vnc-supervisor.conf
[program:x11vnc]
command=x11vnc -forever -usepw -create 
environment=HOME="/nobody",USER="nobody"
user=nobody
