# files/x11vnc-supervisor.conf
[program:x11vnc]
priority=25
directory=/nobody
command=x11vnc -forever -shared -display :1 
environment=HOME="/nobody",USER="nobody"
user=nobody
autostart=true
autorestart=true
stopsignal=QUIT

