# files/xvfb-supervisor.conf
[program:xvfb]
priority=15
directory=/
command=/usr/bin/Xvfb :1 -screen 0 1024x768x16
user=root
autostart=true
autorestart=true
stopsignal=QUIT
